# scripts/supervisor/conf.d/supervisor_word2pdf.conf
[program:word2pdf]

directory=/opt/openoffice4/program
autostart=true
startsecs=5
autorestart=true
startretries=5
command=/opt/openoffice4/program/python t.py