# scripts/supervisor/conf.d/supervisor_newcom_b.conf
[program:newcom_b]

directory=/opt/newcom
autostart=true
startsecs=5
autorestart=true
startretries=5
# environment # 控制环境类型，可选类型为[PRODUCTION,TESTING,]
command=gunicorn -w 9 -b 0.0.0.0:5000 wsgi:app